# packages/deb/DEBIAN/postinst
#!/bin/bash
# postinst script for Nginx Amplify Agent

set -e

if [ "$1" != "configure" ]; then
    exit 0
fi

if [ -z "$2" ]; then
    update-rc.d amplify-agent defaults
    mkdir -p /var/run/amplify-agent
    mkdir -p /var/log/amplify-agent
    touch /var/log/amplify-agent/agent.log
    chown nginx /var/run/amplify-agent
    chown nginx /var/log/amplify-agent
    chown nginx /var/run/amplify-agent /var/log/amplify-agent/agent.log
else
    # Check for an older version of the agent running
    if command -V pgrep > /dev/null 2>&1; then
        agent_pid=`pgrep amplify-agent`
    else
        agent_pid=`ps aux | grep -i '[a]mplify-agent' | awk '{print $2}'`
    fi

    # stop it
    if [ -n "$agent_pid" ]; then
        service amplify-agent stop > /dev/null 2>&1 < /dev/null
    fi

    # Change API URL from 1.0 to 1.1
    sh -c "sed -i.old 's/api_url.*$/api_url = https:\/\/receiver.amplify.nginx.com:443\/1.1/' \
        /etc/amplify-agent/agent.conf"

    # start it
    service amplify-agent start > /dev/null 2>&1 < /dev/null
fi
